# docker-laravel-realtime-chat-app/.docker/supervisor.d/horizon.conf.default
[program:horizon]
process_name=%(program_name)s
command=php /var/www/html/artisan horizon ; command này cũng sẽ chạy ở foreground
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/horizon.log
stopwaitsecs=3600
